# build/code-server.conf
[program:start-script]
autorestart = false
startsecs = 0
user = nobody
command = /usr/local/bin/start.sh
umask = 000

[program:shutdown-script]
autorestart = false
startsecs = 0
user = root
command = /usr/local/bin/shutdown.sh '/usr/lib/code-server'
umask = 000